# src/testdata2/installed_roles_path/os_nova/files/smt.conf
#
# A task to turn off smt
#

description "turn off the smt"

start on runlevel [2345]

task
exec ppc64_cpu --smt=off
