# src/testdata3/scm_repo/playbooks/roles/wsgi-app/templates/supervisor.conf
[program:{{ app_name }}]
command={{ www_dir }}/{{ app_name }}_run.sh
directory={{ app_deploy_dir }}
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-stderr.log
killasgroup=true
stopasgroup=true
